# docker/cron.d/app.cron
* * * * * date >> /var/log/app_cron.log
